# crontab
* * * * * root /code/test.py >> /var/log/task.log 2>&1
# empty line