# crontab
* * * * * root /usr/local/bin/python /code/test.python > /dev/stdout
# empty line