# crontab
* * * * * echo "test" >> /crontest.txt 
# empty line