# crontab
* * * * * python /code/test.python
# empty line