# crontab
* * * * * python /code/test.python > /dev/stdout
# empty line