# crontab
* * * * * /script.sh >> /var/log/script.log
