# app/worker/crontab
* * * * * /app/worker/script.sh >> /var/log/script.log 2>&1
#