# app/worker/crontab
* * * * * /usr/bin/python /app/worker/test.py
#